# assets/content/case_studies/frankenstein.json
{
  "meta": {
    "type": "case_study",
    "slug": "frankenstein",
    "source_path": "security/case_studies/frankenstein_case_study.md",
    "incident_code": "C0001"
  },
  "content": {
    "title": "Frankenstein [C0001]",
    "summary": "Frankenstein was a campaign carried out in early 2019, consisting of highly targeted intrusion activity without destructive payloads or large-scale impact. The attackers pieced together unrelated open-source offensive and administrative tools, rather than using custom malware or zero-days, to achieve their goals. Attribution remains unknown, as nothing in the campaign stood out to identify the threat actors. This case highlights how living-off-the-land and open-source tools can enable stealthy, targeted intrusions, increasing the need for defenders to move beyond signature-based detection.",
    "incident_overview": "Initial access took place via two vectors: first, a trojanized Word document that would download a remote template designed to run arbitrary code using Word Memory Corruption (CVE-2017-11882), then establish persistence via a scheduled task masquerading as \"WinUpdate\". The second vector also used a trojanized Word document, but this time it used macros to launch a Visual Basic script that would first check for malware analysis software and/or signs of being in a virtual machine, then stop execution if such conditions were detected.\n\nFor defense evasion and execution, the attackers abused trusted Windows binaries such as MSBuild, bypassed user controls and application allowlisting, and then deployed a PowerShell Empire post-exploitation framework through an XML file and a stager. The stager also established an encrypted communication channel with a command-and-control (C2) server to leverage PowerShell Empire, and system information was gathered and delivered as well.",
    "adversary_objectives": "The objectives are difficult to determine with certainty, as no destructive actions were taken and stealth appeared to be a priority, resulting in limited impact. The identity of the victims is also unknown, but the campaign was likely specific to enterprise environments heavily dependent on Windows and Active Directory.\n\nThe attackers' actions suggest an interest in espionage or achieving long-term access, consistent with information-gathering objectives. Any additional intent, such as tool testing, remains speculative and cannot be confirmed.",
    "attack_mapping": {
      "tactics": [
        "Command & Control [TA0011]",
        "Collection [TA0009]",
        "Exfiltration [TA0010]",
        "Execution [TA0002]",
        "Defense Evasion [TA0005]",
        "Resource Development [TA0042]",
        "Initial Access [TA0001]",
        "Discovery [TA0007]",
        "Privilege Escalation [TA0004]",
        "Persistence [TA0003]"
      ],
      "techniques": [
        "**Command & Control**: Web Protocols [T1071.001], Symmetric Cryptography [T1573.001], Ingress Tool Transfer [T1105]",
        "**Collection**: Automated Collection [T1119], Data from Local System [T1005], Exfiltration Over C2 Channel [T1041]",
        "**Exfiltration**: Automated Exfiltration [T1020]",
        "**Execution**: PowerShell [T1059.001], Windows Command Shell [T1059.003], Visual Basic [T1059.005], Exploitation for Client Execution [T1203], Scheduled Task [T1053.005], Malicious File [T1204.002], Windows Management Instrumentation [T1047]",
        "**Defense Evasion**: Deobfuscate/Decode Files or Information [T1140], Masquerade Task or Service [T1036.004], Command Obfuscation [T1027.010], Template Injection [T1221], MSBuild [T1127.001]",
        "**Resource Development**: Tool [T1588.002], System Checks [T1497.001]",
        "**Initial Access**: Spearphishing Attachment [T1566.001]",
        "**Discovery**: Process Discovery [T1057], Security Software Discovery [T1518.001], System Information Discovery [T1082], System Network Configuration Discovery [T1016], System Owner/User Discovery [T1033], System Checks [T1497.001]",
        "**Privilege Escalation**: Scheduled Task [T1053.005]",
        "**Persistence**: Scheduled Task [T1053.005]"
      ],
      "procedure_examples": [
        "The campaign used Command and Scripting Interpreter: Visual Basic [T1059.005] and Template Injection [T1221] via a trojanized Word document likely delivered via a Spearphishing Attachment [T1566.001].",
        "The Visual Basic script ran System Checks [T1497.001] to avoid malware analysis tools and virtual machine environments.",
        "A Scheduled Task [T1053.005] was set up to achieve persistence, named \"WinUpdate\" and leveraging the Masquerade Task or Service [T1036.004] technique.",
        "With the use of Web Protocols [T1071.001], a communication channel using Symmetric Cryptography [T1573.001] was established with a C2 server.",
        "The attackers collected a wide variety of information via PowerShell [T1059.001] (Empire tool), such as Process Discovery [T1057], System Information Discovery [T1082], System Network Configuration Discovery [T1016], System Owner/User Discovery [T1033], followed by Automated Exfiltration [T1020] via the established C2 channel."
      ]
    },
    "defensive_analysis": "Phishing emails with malicious attachments were allowed to pass through to organizational email servers, constituting the first defensive failure. Additionally, employees lacked the training or awareness to handle suspicious emails, leading to compromise via malicious Word documents.\n\nOrganizations with mature email filtering, application allowlisting, and script logging likely prevented or disrupted the intrusion.",
    "lessons_learned": [
      "The human element remains a vulnerable aspect of security; ongoing training against phishing tactics is essential, along with robust email filtering and traffic analysis.",
      "Defenders must expect the use of simple and unsophisticated tools, including \"living off the land\" techniques, rather than only preparing for zero-day vulnerabilities or custom malware.",
      "Behavioral detection and monitoring are critical, as signature-based approaches may miss attacks leveraging legitimate or open-source tools."
    ],
    "reflection": "The main value of this case is its simplicity. It shows how ordinary open-source tools and trusted system components can be chained together effectively without relying on custom malware or novel exploits.\n\nIt also reinforces the value of ATT&CK-based reasoning. Mapping the campaign step by step makes it easier to explain why behavioral detection matters, especially when the attacker is relying on tools that can blend into normal system activity.",
    "references": [
      {
        "label": "MITRE ATT&CK: Frankenstein Campaign (C0001)",
        "url": "https://attack.mitre.org/campaigns/C0001/"
      },
      {
        "label": "Attackers stitch together Frankenstein campaign using free tools",
        "url": "https://www.bleepingcomputer.com/news/security/attackers-stitch-together-frankenstein-campaign-using-free-tools/"
      },
      {
        "label": "Frankenstein Campaign",
        "url": "https://blogs.cisco.com/security/talos/frankenstein-campaign"
      },
      {
        "label": "Open Source Frankenstein",
        "url": "https://www.reslam.com/open-source-frankenstein/"
      },
      {
        "label": "MITRE D3FEND Framework",
        "url": "https://d3fend.mitre.org/"
      }
    ],
    "defensive_mapping": "_Scope note: This is a non-exhaustive, conceptual mapping of defensive approaches to observed techniques, not a prescriptive checklist._\n\n- Web Protocols [T1071.001]\n\n  - Hardening: File Encryption\n  - Detection: Outbound Internet Web Traffic, Outbound Internet Network Traffic, Certificate Analysis, Network Traffic Signature Analysis, Network Traffic Community Deviation, Protocol Metadata Anomaly Detection, Client-server Payload Profiling, Per Host Download-Upload Ratio Analysis, Relay Pattern Analysis, Remote Terminal Session Detection, User Geolocation Logon Pattern Analysis, File Integrity Monitoring, File Analysis\n  - Isolation: Outbound Internet Web Traffic, Outbound Internet Network Traffic, Content Modification, Content Quarantine, Network Traffic Filtering, Local File Permissions, Content Filtering, Remote File Access Mediation, Outbound Traffic Filtering\n  - Eviction: File Eviction\n  - Restoration: Restore File\n  - Deception: Decoy File\n\n- Symmetric Cryptography [T1573.001]\n\n  - Detection: Network Traffic Community Deviation, Remote Terminal Session Detection, Network Traffic Signature Analysis, Per Host Download-Upload Ratio Analysis, Relay Pattern Analysis, Protocol Metadata Anomaly Detection, Client-server Payload Profiling, User Geolocation Logon Pattern Analysis\n  - Isolation: Network Traffic Filtering, Outbound Traffic Filtering\n\n- Ingress Tool Transfer [T1105]\n\n  - Detection: Client-server Payload Profiling, Network Traffic Community Deviation, Per Host Download-Upload Ratio Analysis, Protocol Metadata Anomaly Detection, Network Traffic Signature Analysis, Relay Pattern Analysis, Remote Terminal Session Detection, User Geolocation Logon Pattern Analysis\n  - Isolation: Network Traffic Filtering, Outbound Traffic Filtering\n\n- Automated Collection [T1119]\n\n  - Hardening: File Encryption\n  - Detection: File Integrity Monitoring, File Analysis\n  - Isolation: Remote File Access Mediation, Content Quarantine, Local File Permissions, Content Modification, Content Filtering\n  - Eviction: File Eviction\n  - Restoration: Restore File\n  - Deception: Decoy File\n\n- Data from Local System [T1005]\n\n  - Hardening: File Encryption\n  - Detection: File Integrity Monitoring, File Analysis\n  - Isolation: Content Quarantine, Content Modification, Local File Permissions, Content Filtering, Remote File Access Mediation\n  - Eviction: File Eviction\n  - Restoration: Restore File\n  - Deception: Decoy File\n\n- Exfiltration Over C2 Channel [T1041]\n\n  - Hardening: File Encryption\n  - Detection: Network Traffic Community Deviation, Per Host Download-Upload Ratio Analysis, Protocol Metadata Anomaly Detection, Client-server Payload Profiling, Remote Terminal Session Detection, Certificate Analysis, Network Traffic Signature Analysis, File Integrity Monitoring, User Geolocation Logon Pattern Analysis, File Analysis\n  - Isolation: Network Traffic Filtering, Local File Permissions, Content Modification, Content Quarantine, Content Filtering, Remote File Access Mediation\n  - Eviction: File Eviction\n  - Restoration: Restore File\n  - Deception: Decoy File\n\n- Automated Exfiltration [T1020]\n\n  - Detection: Client-server Payload Profiling, Network Traffic Community Deviation, Per Host Download-Upload Ratio Analysis, Protocol Metadata Anomaly Detection, Remote Terminal Session Detection, Network Traffic Signature Analysis, User Geolocation Logon Pattern Analysis\n  - Isolation: Network Traffic Filtering\n\n- Command and Scripting Interpreter [T1059]\n\n  - Hardening: File Encryption\n  - Detection: Emulated File Analysis, Dynamic Analysis, File Integrity Monitoring, File Analysis\n  - Isolation: Content Quarantine, Content Modification, Executable Denylisting, Executable Allowlisting, Content Filtering, Remote File Access Mediation, Local File Permissions\n  - Eviction: File Eviction\n  - Restoration: Restore File\n  - Deception: Decoy File\n\n- Exploitation for Client Execution [T1203]\n\n  - Hardening: Stack Frame Canary Validation, Process Segment Execution Prevention, Segment Address Offset Randomization\n  - Detection: Process Code Segment Verification, Shadow Stack Comparisons, Memory Boundary Tracking\n\n- Scheduled Task [T1053.005]\n\n  - Detection: Process Spawn Analysis, System Call Analysis, Process Self-Modification Detection, System Daemon Monitoring, Scheduled Job Analysis\n  - Isolation: Application-based Process Isolation, Hardware-based Process Isolation, Kernel-based Process Isolation, Executable Allowlisting, Executable Denylisting, System Call Filtering\n  - Eviction: Process Termination, Process Suspension, Host Shutdown, Host Reboot\n\n- Malicious File [T1204.002]\n\n  - Hardening: File Encryption\n  - Detection: Emulated File Analysis, Dynamic Analysis, Remote Terminal Session Detection, Network Traffic Signature Analysis, Per Host Download-Upload Ratio Analysis, Relay Pattern Analysis, Network Traffic Community Deviation, Client-server Payload Profiling, Protocol Metadata Anomaly Detection, Homoglyph Detection, Identifier Activity Analysis, User Geolocation Logon Pattern Analysis, URL Analysis, File Integrity Monitoring, File Analysis, URL Reputation Analysis\n  - Isolation: Executable Denylisting, Executable Allowlisting, Network Traffic Filtering, Local File Permissions, Content Quarantine, Content Modification, Content Filtering, Remote File Access Mediation, Outbound Traffic Filtering\n  - Eviction: File Eviction\n  - Restoration: Restore File\n  - Deception: Decoy File\n\n- Windows Management Instrumentation [T1047]\n\n  - Detection: Protocol Metadata Anomaly Detection, Remote Terminal Session Detection, Network Traffic Community Deviation, System Call Analysis, Process Spawn Analysis, Administrative Network Activity Analysis, Network Traffic Signature Analysis, Per Host Download-Upload Ratio Analysis, Client-server Payload Profiling, User Geolocation Logon Pattern Analysis, Connection Attempt Analysis\n  - Isolation: Executable Denylisting, Executable Allowlisting, System Call Filtering, Network Traffic Filtering, Hardware-based Process Isolation\n\n- Deobfuscate/Decode Files or Information [T1140]\n\n  - Hardening: File Encryption\n  - Detection: Dynamic Analysis, Emulated File Analysis, Process Spawn Analysis, System Call Analysis, File Integrity Monitoring, File Analysis\n  - Isolation: Local File Permissions, Executable Denylisting, Hardware-based Process Isolation, Executable Allowlisting, Content Modification, Content Quarantine, System Call Filtering, Content Filtering, Remote File Access Mediation\n  - Eviction: File Eviction\n  - Restoration: Restore File\n  - Deception: Decoy File\n\n- Masquerade Task or Service [T1036.004]\n\n  - Detection: Scheduled Job Analysis\n\n- MSBuild [T1127.001]\n\n  - Hardening: File Encryption, Software Update\n  - Detection: File Analysis, File Integrity Monitoring\n  - Isolation: Local File Permissions, Content Modification, Content Quarantine, Content Filtering, Remote File Access Mediation\n  - Eviction: File Eviction\n  - Modeling: Asset Vulnerability Enumeration, Software Inventory\n  - Restoration: Restore File, Restore Software\n  - Deception: Decoy File\n\n- Spearphishing Attachment [T1566.001]\n\n  - Hardening: File Encryption\n  - Detection: Dynamic Analysis, Emulated File Analysis, Client-server Payload Profiling, Network Traffic Community Deviation, Per Host Download-Upload Ratio Analysis, Protocol Metadata Anomaly Detection, Inbound Session Volume Analysis, Network Traffic Signature Analysis, Remote Terminal Session Detection, File Integrity Monitoring, Sender MTA Reputation Analysis, Sender Reputation Analysis, User Geolocation Logon Pattern Analysis, Homoglyph Detection, File Analysis\n  - Isolation: Content Quarantine, Content Modification, Local File Permissions, Network Traffic Filtering, Content Filtering, Remote File Access Mediation, Email FIltering, Inbound Traffic Filtering\n  - Eviction: File Eviction, Email Removal\n  - Modeling: Data Inventory\n  - Restoration: Restore File, Restore Email\n  - Deception: Decoy File\n\n- Process Discovery [T1057]\n\n  - Detection: System Call Analysis, Process Spawn Analysis\n  - Isolation: Hardware-based Process Isolation, Executable Allowlisting, Executable Denylisting, System Call Filtering\n\n- Security Software Discovery [T1518.001]\n\n  - Detection: System Call Analysis\n  - Isolation: Content Quarantine, System Call Filtering\n  - Modeling: Configuration Inventory\n  - Restoration: Restore Configuration\n\n- System Information Discovery [T1082]\n\n  - Detection: System Call Analysis, Process Spawn Analysis\n  - Isolation: Executable Allowlisting, Executable Denylisting, Hardware-based Process Isolation, System Call Filtering\n  - Deception: Decoy Environment\n\n- System Network Configuration Discovery [T1016]\n\n  - Hardening: File Encryption\n  - Detection: System Call Analysis, Process Spawn Analysis, Dynamic Analysis, Emulated File Analysis, File Integrity Monitoring, File Analysis\n  - Isolation: Content Quarantine, Content Modification, Executable Denylisting, Hardware-based Process Isolation, Executable Allowlisting, Local File Permissions, System Call Filtering, Content Filtering, Remote File Access Mediation\n  - Eviction: File Eviction\n  - Restoration: Restore File\n  - Deception: Decoy File\n\n- System Owner/User Discovery [T1033]\n\n  - Hardening: Process Segment Execution Prevention, Segment Address Offset Randomization, File Encryption\n  - Detection: Process Self-Modification Detection, System Call Analysis, Process Spawn Analysis, File Integrity Monitoring, File Analysis, Process Lineage Analysis\n  - Isolation: Content Quarantine, Content Modification, Executable Allowlisting, Executable Denylisting, Hardware-based Process Isolation, Application-based Process Isolation, Kernel-based Process Isolation, Domain Trust Policy, Local File Permissions, System Call Filtering, Content Filtering, Remote File Access Mediation, Web Session Access Mediation\n  - Eviction: Process Termination, Process Suspension, Host Shutdown, File Eviction, Host Reboot\n  - Restoration: Restore File, Restore Database\n  - Deception: Decoy File",
    "atlas_mapping": {
      "summary": "ATLAS does not apply to this incident."
    },
    "indicators": {
      "items": [
        "Phishing emails",
        "Trojanized Word documents",
        "Scheduled Task masquerading as \"WinUpdate\"",
        "Encrypted outbound web traffic consistent with PowerShell Empire C2 activity."
      ],
      "summary": "Public reporting did not disclose concrete hashes, IP addresses, or domains for this campaign. Behavioral indicators include:"
    }
  }
}
